# systemd/xonotic-all.target
[Unit]
Discription=All Xonotic Servers
Requires= # xonotic@dm.service xonotic@duel.service etc

[Install]
WantedBy=multi-user.target
